# license-server/deploy/license-server-backup.service
[Unit]
Description=License server SQLite backup
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/opt/license-server
ExecStart=/usr/bin/npm run backup
User=licenseserver
